# deploy/eve-trade-hub-sde-sync.service
[Unit]
Description=EVE trade hub analyzer — weekly SDE sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=eve-hub
Group=eve-hub
WorkingDirectory=/var/lib/eve-trade-hub-analyzer
EnvironmentFile=/etc/eve-trade-hub-analyzer/env
ExecStart=/usr/local/bin/sde-sync
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
